# data/motds/toriel/update-motd.d/25-baller-motd
#!/usr/bin/env bash
#first part will create 'center info'
#final part prints total outcome
DIR="$( cd "$( dirname "${BASH_SOURCE[0]}" )" && pwd )"
cd $DIR
cd ./motd_art/
cat toriel.txt
